# mysite/supervisord/supervisor-app.conf
[supervisord]
nodaemon=true 


[supervisorctl]

[program:app-uwsgi]
command = /usr/local/bin/uwsgi --ini /home/docker/code/uwsgi/uwsgi.ini

[program:nginx-app]
command = /usr/sbin/nginx